# scripts/logrotate.cron
0 0 * * * Spring-Hello-Log4j/scripts/logrotate.sh
